# datamart/stg/AMK_PROD.crontab
00 23 * * 1-4 source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_daily_sched_prod.xml 
00 20 * * 5 source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_daily_sched_prod.xml
59 23 * * * source /amk/projects/default/amkmietl/bin/crontab_backup.sh 1>~/crontab_backup.log 2>&1
01 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --file_handler_name "Data Retention AMO-530"
05 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amk/projects/default/amkmietl/conf/AMK_PROD.json --file_handler "AMK-504 compressed file clean up"
05 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amk/projects/default/amkmietl/conf/AMK_PROD.json --file_handler "AMK-504 log clean up"
00 13 * * 0 source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_weekly_sched_prod.xml
00 19 * * 6 source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/sat_maintenance.xml
*/3 * * * * /usr/bin/python /amk/projects/default/amkmietl/ssod-data-quality-package/svn/auto_deploy.py -l true -e 'ssoamkdev@sas.com' -p '/amk/projects/default/amkmietl/' -d true >> /amk/projects/default/amkmietl/logs/svn_auto_deploy.log 2>&1
*/15 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status_AMK.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status_AMK.log 2>&1
  30 21 * * 0  /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.log 2>&1
30 2 * * 5 /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh > /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh.log
13 20 23 08 * source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_weekly_sched_prod.xml --skip_dependencies
40 20 17 01 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
02 21 14 02 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
04 21 20 03 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
46 20 29 05 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
48 20 19 06 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
38 20 03 07 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
41 20 10 07 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
21 21 02 08 * source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_weekly_sched_prod.xml
02 22 15 08 * cd /amk/projects/default/amkmietl; /amk/projects/default/amkmietl/bin/sso_batch.sh -c /amk/projects/default/amkmietl/conf/amk_config.txt -f -s /amk/projects/default/amkmietl/conf/amk_sync_flow1_mod2prod.xml
48 20 13 09 * source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_weekly_sched_prod.xml --skip_dependencies
11 21 20 09 * source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_weekly_sched_prod.xml --resume --skip_dependencies
