# datamart/stg/AMK_PROD-test.crontab
00 23 * * 1-4 source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_daily_sched_prod.xml 
00 20 * * 5 source /amk/projects/default/amkmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --schedule=/amk/projects/default/amkmietl/conf/amk_daily_sched_prod.xml
01 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file=/amk/projects/default/amkmietl/conf/AMK_PROD.json --file_handler_name "Data Retention AMO-530"
05 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amk/projects/default/amkmietl/conf/AMK_PROD.json --file_handler "AMK-504 compressed file clean up"
05 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amk/projects/default/amkmietl/conf/AMK_PROD.json --file_handler "AMK-504 log clean up"
30 21 * * 0  /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.log 2>&1
30 2 * * 5 /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh > /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh.log
