# datamart/stg/AMU_PROD.crontab
00 20 * * 1-5 source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/amu_daily_sched_prod.xml
00 12 * * 0 source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/amu_weekly_sched_prod.xml
45 18 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amu/projects/default/amumietl/conf/AMU_PRD.json --file_handler "Compress_backup_directory_AMU-580"
45 18 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amu/projects/default/amumietl/conf/AMU_PRD.json --file_handler "AMU_delete_DM_tar_file"
45 18 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amu/projects/default/amumietl/conf/AMU_PRD.json --file_handler "AMU_delete_DM_log_file"
00 19 * * 6 source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/sat_maintenance.xml
*/3 * * * * /usr/bin/python /amu/projects/default/amumietl/ssod-data-quality-package/svn/auto_deploy.py -l true -e 'ssoamudev@sas.com' -p '/amu/projects/default/amumietl/' -d true >> /amu/projects/default/amumietl/logs/svn_auto_deploy.log 2>&1
01 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amu/projects/default/amumietl/conf/AMU_PRD.json --file_handler_name "Data Retention AMO-530"
59 23 * * * source /amu/projects/default/amumietl/bin/crontab_backup.sh 1>~/crontab_backup.log 2>&1
*/15 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status_AMU.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status_AMU.log 2>&1
*/5 * * * * /home/scnzqz/SEM/crontab_export.sh
30 21 * * 0  /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.log 2>&1
30 2 * * 5 /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh > /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh.log
55 20 7 2 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
38 20 18 03 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow3_prod2psp.xml
51 20 08 05 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
50 20 29 05 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow2_prod2mod.xml
20 08 04 05 * source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/amu_weekly_sched_prod.xml
47 18 23 08 * source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/amu_weekly_sched_prod.xml --skip_dependencies
58 20 30 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -s /amu/projects/default/amumietl/conf/amu_weekly_sched_prod.xml
44 20 19 06 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
42 20 26 06 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
47 20 10 07 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
59 21 10 07 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow3_prod2psp.xml
20 21 07 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
33 20 20 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow3_prod2psp.xml
29 23 20 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -s /amu/projects/default/amumietl/conf/amu_sync_flow3_prod2psp.xml
08 21 21 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
14 22 21 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow3_prod2psp.xml 
38 20 28 08 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
19 21 20 09 * source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/amu_weekly_sched_prod.xml --resume --skip_dependencies
14 16 06 09 * source /amu/projects/default/amumietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amu/projects/default/amumietl/conf/AMU_PRD.json --schedule=/amu/projects/default/amumietl/conf/amu_weekly_sched_prod.xml --resume --skip_dependencies
33 20 11 09 * cd /amu/projects/default/amumietl; /amu/projects/default/amumietl/bin/sso_batch.sh -c /amu/projects/default/amumietl/conf/amu_config.txt -f -s /amu/projects/default/amumietl/conf/amu_sync_flow1_mod2prod.xml
