# datamart/stg/AMV_PROD.crontab
45 18 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amv/projects/default/amvmietl/conf/AMV_PRD.json --file_handler "Compress_backup_directory_AMV-474"
45 18 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amv/projects/default/amvmietl/conf/AMV_PRD.json --file_handler "AMV_delete_DM_tar_file"
45 18 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amv/projects/default/amvmietl/conf/AMV_PRD.json --file_handler "AMV_delete_DM_log_file"
00 19 * * 6 source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/sat_maintenance.xml
*/3 * * * * /usr/bin/python /amv/projects/default/amvmietl/ssod-data-quality-package/svn/auto_deploy.py -l true -e 'ssoamvdev@sas.com' -p '/amv/projects/default/amvmietl/' -d true >> /amv/projects/default/amvmietl/logs/svn_auto_deploy.log 2>&1
01 07 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /amv/projects/default/amvmietl/conf/AMV_PRD.json --file_handler_name "Data Retention AMO-530"
59 23 * * * source /amv/projects/default/amvmietl/bin/crontab_backup.sh 1>~/crontab_backup.log 2>&1
00 22 * * 1-5 source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_daily_sched_prod.xml
00 12 * * 0 source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml
*/15 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status_AMV.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status_AMV.log 2>&1
30 21 * * 0  /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/src/custom/runtime_report/runtime_report.log 2>&1
30 2 * * 5 /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh > /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh.log
38 22 25 09 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
46 21 30 01 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow3_prod2psp.xml
00 08 17 02 * source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --resume  --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml
38 13 17 02 * source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml
10 22 23 02 * source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml
33 21 29 05 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
16 20 10 05 * source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml --skip_dependencies
46 21 19 06 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
32 21 03 07 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
31 18 02 08 * source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml
16 22 15 08 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
41 20 16 08 * source /amv/projects/default/amvmietl/LoadMgr/etc/.LoadMgr_profile; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/amv/projects/default/amvmietl/conf/AMV_PRD.json --schedule=/amv/projects/default/amvmietl/conf/amv_weekly_sched_prod.xml
36 22 28 08 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
31 22 11 09 * cd /amv/projects/default/amvmietl; /amv/projects/default/amvmietl/bin/sso_batch.sh -c /amv/projects/default/amvmietl/conf/amv_config.txt -f -s /amv/projects/default/amvmietl/conf/amv_sync_flow1_mod2prod.xml
