# datamart/stg/ANF_PROD.crontab
*/5 * * * * /anf/projects/default/anfmi/bin/sso_run_error_check.py -c /anf/projects/default/anfmi/conf/anf_config.txt --look-back=10 --log-level=DEBUG --locate-zabbix-args 
59 23 * * * source /anf/projects/default/anfmi/bin/crontab_backup.sh >~/crontab_backup.log 2>&1
*/3 * * * * /usr/bin/python /anf/projects/default/anfmi/ssod-data-quality-package/svn/auto_deploy.py -l true -e 'ssoanfdev@sas.com' -p '/anf/projects/default/anfmi/' -d true >> /anf/projects/default/anfmi/logs/svn_auto_deploy.log 2>&1
30 21 * * 0-4 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_daily_sched_prod.log 
30 21 * * 5 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_friday_night_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_friday_night_sched_prod.log 
30 21 * * 5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler "ANF-2512_Compress_backup_directory"
31 21 * * 5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler "ANF_DM_backup_file_clean_up"
33 21 * * 5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler "ANF_DM_backup_log_file_clean_up"
0 5 * * 1-6 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_results_sla_check_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_results_sla_check_prod.log 
30 7 * * 1-6 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_batch_complete_sla_check_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_batch_complete_sla_check_prod.log 
30 14 * * 0 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_results_sla_check_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_results_sla_check_prod.log 
30 02 * * 0 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_weekly_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_weekly_sched_prod.log 
0 19 * * 0-5 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_profile_export_and_purge_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_profile_export_and_purge_prod.log 
0 19 * * 6 source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_data_extract_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_data_extract_prod.log
0 15 * * 0-6 /usr/bin/python /anf/projects/default/anfmi/bin/sso_clean_incoming_dir.py >> /anf/projects/default/anfmi/logs/sso_clean_incoming_dir.log 2>&1
0 15 * * 0-6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler_name ANF_Transport_1M__clean_up
1 15 * * 0-6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler_name ANF_Transport_1W__clean_up
2 15 * * 0-6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler_name ANF_Transport_1W_Directory_clean_up
3 15 * * 0-6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_handler_name ANF_Transport_1M_Directory_clean_up
0 18 * * * /usr/bin/python /anf/projects/default/anfmi/bin/manage_files.py 20  10  /sso/transport/archive >> /anf/projects/default/anfmi/logs/sso_manage_files.log 2>&1
0 19 * * * /usr/bin/python /anf/projects/default/anfmi/bin/manage_files.py 30  60 /sso/sfw/LoadMgr/logs/backup >> /anf/projects/default/anfmi/logs/sso_manage_files.log 2>&1
* * * * * source /sso/sfw/ghusps-toolkit/toolkit_venv/bin/activate; cd /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers; python ams_route_files.py --config_file=/anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_route_name=anf_sftp_realtime_pull_prod >> /sso/sfw/LoadMgr/logs/stb_requests/anf_sftp_realtime_pull_prod.log 2>&1
* * * * * source /sso/sfw/ghusps-toolkit/toolkit_venv/bin/activate; cd /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers; python ams_route_files.py --config_file=/anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_route_name=anf_sftp_realtime_push_prod >> /sso/sfw/LoadMgr/logs/stb_requests/anf_sftp_realtime_push_prod.log 2>&1
*/10 * * * * source /sso/sfw/ghusps-toolkit/toolkit_venv/bin/activate; cd /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers; python ams_route_files.py --config_file=/anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --file_route_name=anf_needs_processing_retry >> /sso/sfw/LoadMgr/logs/stb_requests/anf_needs_processing_retry.log 2>&1
30 21 * * 0 /sso/sfw/LoadMgr/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/custom/runtime_report/runtime_report.log 2>&1
30 08 * * * /sso/sfw/LoadMgr/custom/runtime_report_daily/runtime_report.sh > /sso/sfw/LoadMgr/custom/runtime_report_daily/runtime_report.log 2>&1
30 2 * * 5 /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh > /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh.log
1,16,31,46 19-23,0-9 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status.log 2>&1
30 8 * * * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /anf/projects/default/anfmi/bin/search4_invalid_http_req.py
16 2 22 7 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --resume --schedule /anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_daily_sched_prod.log
57 11 19 07 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule=/anf/projects/default/anfmi/conf/anf_data_extract_prod.xml --adhoc_schedule=/anf/projects/default/anfmi/conf/anf_complete_backup.xml >> /anf/projects/default/anfmi/logs/anf_complete_backup.xml_`date +\%Y\%m\%d_\%H\%M\%S`.log
23 3 22 7 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_daily_sched_prod.log 
30 0 25 7 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_friday_night_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_friday_night_sched_prod.log
30 13 04 08 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml --adhoc_schedule=/anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml_anf_2765 2>&1 >> /anf/projects/default/anfmi/logs/anf_daily_sched_prod.xml_anf_2765
45 04 06 08 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --resume --skip_dependencies --schedule /anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_daily_sched_prod.log
39 21 23 09 * source /sso/sfw/LoadMgr/etc/environment_config.sh; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /anf/projects/default/anfmi/conf/ams_toolkit_config_2.json --schedule /anf/projects/default/anfmi/conf/anf_daily_sched_prod.xml 2>&1 >> /anf/projects/default/anfmi/logs/anf_daily_sched_prod.log
