# datamart/stg/ASK_PROD.crontab
30 17 * * 0 cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --schedule /ask/projects/default/askmi/conf/ask_weekly_sched_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
*/5 * * * * /ask/projects/default/askmi/bin/sso_run_error_check.py -c /ask/projects/default/askmi/conf/ask_config.txt --look-back=10 --log-level=DEBUG --locate-zabbix-args 2>&1 >> /ask/projects/default/askmi/logs/sso_run_error_check.log
59 23 * * * source /ask/projects/default/askmi/bin/crontab_backup.sh >~/crontab_backup.log 2>&1
00 22 * * 6 /ask/projects/default/askmi/bin/purge_dm_backup_PROD2.sh > /ask/projects/default/askmi/logs/purge_dm_backup_PROD2.log
01 22 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python  /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --file_handler_name "ASK_DM_backup_file_clean_up"
02 22 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python  /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --file_handler_name "ASK_DM_backup_log_file_clean_up"
03 22 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python  /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --file_handler_name "ASK_DM_backup_err_file_clean_up"
00 08 * * 0 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python  /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --file_handler_name "ASM-213_Compress_backup_directory"  >> /ask/warehouse/default/logs/maintenance/ASK_COMPRESS_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
05 23 * * 6 cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --schedule /ask/projects/default/askmi/conf/ask_saturday_maint_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
15,45 05-19 * * 1-5 cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --schedule /ask/projects/default/askmi/conf/b8_intraday_sched_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
0 21 * *  1-5 cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --schedule /ask/projects/default/askmi/conf/b8_nightly_sched_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
01,16,31,46 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh.log 2>&1
49 21 06 09 * cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --resume --skip_dependencies  --schedule /ask/projects/default/askmi/conf/ask_weekly_sched_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
4 22 21 11 * cd /ask/projects/default/askmi; /ask/projects/default/askmi/bin/sso_batch.sh -c /ask/projects/default/askmi/conf/ask_config.txt -f -s /ask/projects/default/askmi/conf/ask_complete_backup_compress.xml
10 01 29 11 * cd /ask/projects/default/askmi; /ask/projects/default/askmi/bin/sso_batch.sh -c /ask/projects/default/askmi/conf/ask_config.txt -f -s /ask/projects/default/askmi/conf/ask_complete_backup_compress.xml
16 21 08 03 * cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --resume --schedule /ask/projects/default/askmi/conf/ask_weekly_sched_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
18 20 21 06 * cd /ask/projects/default/askmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /ask/projects/default/askmi/conf/askprodconfig.json --resume   --skip_dependencies    --schedule /ask/projects/default/askmi/conf/ask_weekly_sched_prod.xml >> /ask/warehouse/default/logs/ASK_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
53 13 24 07 * cd /ask/projects/default/askmi; /ask/projects/default/askmi/bin/sso_batch.sh -c /ask/projects/default/askmi/conf/ask_config.txt -f -s /ask/projects/default/askmi/conf/ask_saturday_maint.xml
