# datamart/stg/ASM_PROD.crontab
*/5 * * * * /asm/projects/default/asmmi/bin/sso_run_error_check.py -c /asm/projects/default/asmmi/conf/asm_config.txt --look-back=10 --log-level=DEBUG --locate-zabbix-args 2>&1 >> /asm/projects/default/asmmi/logs/sso_run_error_check.log
59 23 * * * source /asm/projects/default/asmmi/bin/crontab_backup.sh >~/crontab_backup.log 2>&1
30 14 * * 0 cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/asm_weekly_sched_prod.xml >> /asm/warehouse/default/logs/ASM_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
0 21 * * 1-5 cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/b8_nightly_sched_prod.xml >> /asm/warehouse/default/logs/ASM_NIGHTLY_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
00 22 * * 6 /asm/projects/default/asmmi/bin/purge_dm_backup_PROD2.sh > /asm/projects/default/asmmi/logs/purge_dm_backup_PROD2.log
50 23 * * 6 cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/asm_saturday_maint_prod.xml >> /asm/warehouse/default/logs/ASM_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
51 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file=/asm/projects/default/asmmi/conf/asmprodconfig.json --file_handler "ASM-193 backup file clean up"
51 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file=/asm/projects/default/asmmi/conf/asmprodconfig.json --file_handler "ASM-193 clean up of log files related DM backups"
00 08 * * 0 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file=/asm/projects/default/asmmi/conf/asmprodconfig.json --file_handler "ASM-213_Compress_backup_directory"   >> /asm/warehouse/default/logs/maintenance/ASM_COMPRESS_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
01,16,31,46 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh.log 2>&1
37 22 27 11 * cd /asm/projects/default/asmmi; /asm/projects/default/asmmi/bin/sso_batch.sh -c /asm/projects/default/asmmi/conf/asm_config.txt -f -s /asm/projects/default/asmmi/conf/asm_complete_backup_compress.xml
11 22 16 12 * cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/b8_nightly_sched_prod.xml --resume >> /asm/warehouse/default/logs/ASM_NIGHTLY_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
47 21 06 05 * cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/b8_nightly_sched_prod.xml >> /asm/warehouse/default/logs/ASM_NIGHTLY_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
20 20 24 05 * cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/asm_weekly_sched_prod.xml >> /asm/warehouse/default/logs/ASM_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
55 13 24 07 * cd /asm/projects/default/asmmi; /asm/projects/default/asmmi/bin/sso_batch.sh -c /asm/projects/default/asmmi/conf/asm_config.txt -f -s /asm/projects/default/asmmi/conf/asm_saturday_maint_prod.xml
45 21 26 07 * cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/asm_weekly_sched_prod.xml --resume --skip_dependencies >> /asm/warehouse/default/logs/ASM_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
28 22 09 08 * cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/asm_weekly_sched_prod.xml --resume --skip_dependencies >> /asm/warehouse/default/logs/ASM_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
54 21 16 08 * cd /asm/projects/default/asmmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asm/projects/default/asmmi/conf/asmprodconfig.json --schedule /asm/projects/default/asmmi/conf/asm_weekly_sched_prod.xml --resume --skip_dependencies >> /asm/warehouse/default/logs/ASM_SPO_34_PROD_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
