# datamart/stg/ASU_PROD.crontab
59 23 * * * source /asu/projects/default/asumi/bin/crontab_backup.sh >~/crontab_backup.log 2>&1
30 21 * * 0 cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_weekly_sched_prod.xml
*/5 * * * * /asu/projects/default/asumi/bin/sso_run_error_check.py -c /asu/projects/default/asumi/conf/asu_config.txt --look-back=10 --log-level=DEBUG --locate-zabbix-args 2>&1 >> /asu/projects/default/asumi/logs/sso_run_error_check.log
15,45 05-19 * * 1-5 cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/b8_intraday_sched_prod.xml
0 21 * * 1-5 cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/b8_nightly_sched_prod.xml
00 08 * * 0 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --file_handler "ASU-227_Compress_backup_directory"
01 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --file_handler "ASU_delete_DM_tar_file"
01 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --file_handler "ASU_delete_DM_log_file"
50 23 * * 6 cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_saturday_maint_prod.xml
01,16,31,46 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh.log 2>&1
6 1 21 09 * cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_weekly_sched_prod.xml --resume --skip_dependencies
55 13 24 07 * cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_saturday_maint_prod.xml
12 16 08 03 * /asu/projects/default/asumi/LoadMgr/src/so_dm_backup_adhoc/compress_purge_dm_backup.sh > /tmp/compress.log
31 01 21 09 * cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_weekly_sched_prod.xml --resume --skip_dependencies
56 13 22 09 * cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_weekly_sched_prod.xml --resume --skip_dependencies
32 13 24 09 * cd /asu/projects/default/asumi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asu/projects/default/asumi/conf/ASU_SPO_34_PROD.json --schedule /asu/projects/default/asumi/conf/asu_weekly_sched_prod.xml --resume --skip_dependencies
