# datamart/stg/ASV_PROD.crontab
*/5 * * * * /asv/projects/default/asvmi/bin/sso_run_error_check.py -c /asv/projects/default/asvmi/conf/asv_config.txt --look-back=10 --log-level=DEBUG --locate-zabbix-args 2>&1 >> /asv/projects/default/asvmi/logs/sso_run_error_check.log
59 23 * * * source /asv/projects/default/asvmi/bin/crontab_backup.sh >~/crontab_backup.log 2>&1
59 23 * * * cd /sso/sfw/LoadMgr/logs; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --file_handler Compress_HLM_Logs  --config_file=/asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json
0 15 * * 0 cd /asv/projects/default/asvmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --schedule /asv/projects/default/asvmi/conf/asv_weekly_sched_prod.xml
*/5 * * * 0 /asv/projects/default/asvmi/bin/warehouse_monitor.sh >> /asv/projects/default/logs/warehouse_monitor.log 
15,45 07-19 * * 1-5 cd /asv/projects/default/asvmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --schedule /asv/projects/default/asvmi/conf/b8_intraday_sched_prod.xml
0 21 * *  1-5 cd /asv/projects/default/asvmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --schedule /asv/projects/default/asvmi/conf/b8_nightly_sched_prod.xml
00 22 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --file_handler "269_delete_backup_directory"
06 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --file_handler "ASV_Compress_backup_directory"
07 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --file_handler "ASV_delete_DM_tar_file"
08 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --file_handler "ASV_delete_DM_log_file"
05 23 * * 6 cd /asv/projects/default/asvmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --schedule /asv/projects/default/asvmi/conf/asv_saturday_maint_prod.xml
01,16,31,46 19-23,0-8 * * * /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_status.sh.log 2>&1
00 22 28 06 * cd /asv/projects/default/asvmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --schedule /asv/projects/default/asvmi/conf/asv_weekly_sched_prod.xml --resume --skip_dependencies 
56 13 24 07 * cd /asv/projects/default/asvmi; /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /asv/projects/default/asvmi/conf/ASV_SPO_34_PROD.json --schedule /asv/projects/default/asvmi/conf/asv_saturday_maint_prod.xml
