# datamart/stg/CHO_PROD.crontab
55 * * * * /cho/projects/ci/chorun/bin/launch_checksignals.sh
2,17,32,47 * * * * /cho/projects/ci/chorun/bin/sso_zabbix_monitor_logs.sh
51 9 * * * cd /cho/warehouse/data/temp30day;find . -type f -mtime +30 -exec rm -f {} \;  >> /cho/warehouse/data/temp30day/cleanup.log 2>&1
00 00 1 3,6,9,12 * find /cho/projects/default/chorun/logs -mindepth 1 -mtime +190 | xargs rm -rf >> /cho/projects/default/chorun/logs/cleanup.log 2>&1
*/10 20-23,0-9 * * * sh /cho/projects/ci/chorun/bin/check_schedule_status_CHO.sh > /tmp/check_schedule_status_CHO.sh.log 2>&1
00  20 * * MON,TUE,WED,THU,SAT /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/cho_etl_schedule_wkns.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
30  23  * * MON,TUE,WED,THU,SAT /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/individuating_schedule_wkns.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
00 15 * * FRI  /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/cho_etl_schedule_friday.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
30 18  * * FRI /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/individuating_schedule_friday.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
15  13 * * SUN  /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/cho_etl_schedule_sunday.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
00  17  * * SUN  /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/individuating_schedule_sunday.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
00  19  * * SUN  /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/fiscal_week_schedule_w_sunday.xml  >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
00  14  * * MON,TUE,WED,THU,FRI,SAT  /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/extract_ready_early_warning.xml  >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
00  12  * * SUN  /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/extract_ready_early_warning.xml  >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
27  17  7  07 *   /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/extract_ready_early_warning.xml  >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
52 23 08 07 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --resume --schedule=/cho/projects/ci/chorun/conf/cho_etl_schedule_wkns.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
55 23 09 07 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/cho/projects/ci/chorun/conf/config_prod.json --schedule=/cho/projects/ci/chorun/conf/individuating_schedule_wkns.xml >> /cho/projects/ci/chorun/logs/toolkit/toolkit_`date +\%Y\%m\%d_\%H\%M\%S`.log 2>&1
