# datamart/stg/NSO_PROD.crontab
00 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_complete_backup_compress.xml
00 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --file_handler "NSO-546 DM backup file clean up"
00 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --file_handler "NSO clean up of log files related DM backups"
00 23 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --file_handler "NSO clean up of err files related DM backups"
00 12 * * 0 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
00 18 * * 5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/geo_synch.xml
00 12,17 * * 1-5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_message_export.xml
*/5 * * * * cd /nso/projects/default/scripts; /nso/projects/default/scripts/nsomi.5min > /nso/projects/default/logs/nsomi.5min.log 2>&1
30 21 * * 0 /sso/sfw/LoadMgr/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/custom/runtime_report/runtime_report.log 2>&1
*/15 0-8,19-23 * * * /sso/sfw/LoadMgr/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/custom/check_status/check_status.log 2>&1
* * * * * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_route_files.py --config_file=/nso/projects/default/nsomi/conf/nsoprd.json --file_route_name=SftpPushToAWSordcmt
30 23 * * * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file=/nso/projects/default/nsomi/conf/nsoprd.json --file_handler=purge_archive
41 21 02 02 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --resume --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
00 18 07 02 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
00 18 13 02 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
33 20 15 03 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --resume --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
00 18 18 03 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --resume --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
49 19 18 03 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
30 16 06 04 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
15 16 07 04 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
15 19 12 04 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
00 19 14 04 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
16 13 18 05 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_message_export.xml
00 19 06 05 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_complete_backup_compress.xml
33 18 26 07 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file /nso/projects/default/nsomi/conf/nsoprd.json --schedule /nso/projects/default/nsomi/conf/nso_weekly.xml
