# static/temp/JCP_PRD.crontab
30 17 * * * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /jcp/projects/hpmdo/conf/JCP_MI_PROD.json --file_handler "ErrBackupRentention"
30 17 * * * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /jcp/projects/hpmdo/conf/JCP_MI_PROD.json --file_handler "LogBackupRentention"
30 17 * * * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_file_handler.py --config_file /jcp/projects/hpmdo/conf/JCP_MI_PROD.json --file_handler "GzBackupRentention"
30 08 * * 0 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sun.xml
30 08 * * 0 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sun.xml
59 12 * * 1-5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_week_night_reports.xml
30 09 * * 1-5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_week_night_reports.xml
00 18 * * 1 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_mon.xml
00 18 * * 2 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_tues.xml
30 19 * * 3 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_wed.xml
30 19 * * 4 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_thurs.xml
59 15 * * 5 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_fri_backups_psp2prod_sync.xml
00 20 * * 6 /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml
00 23 * * 6 /jcp/projects/default/bin/purge_LM_backup_logs.sh
*/5 * * * * /sso/sfw/LoadMgr/cronshells/jcpetlusr.5min > /sso/sfw/LoadMgr/cronshells/jcpetlusr.5min.log 2>&1
30 21 * * 0  /sso/sfw/LoadMgr/custom/runtime_report/runtime_report.sh > /sso/sfw/LoadMgr/custom/runtime_report/runtime_report.log 2>&1
*/15 * * * * /sso/sfw/LoadMgr/custom/process_time/process_time.sh > /sso/sfw/LoadMgr/custom/process_time/process_time.15min.log 2>&1
01,16,31,46 20-23,0-8 * * * /sso/sfw/LoadMgr/custom/check_status/check_status.sh > /sso/sfw/LoadMgr/custom/check_status/check_status.15min.log 2>&1
15,45 0-5 * * 1-5 sh /sso/sfw/LoadMgr/src/custom/check_status/check_geo_prod_agg.sh > /sso/sfw/LoadMgr/src/custom/check_status/check_geo_prod_agg.log 2>&1
45 01 * * 5 python /sso/sfw/LoadMgr/src/custom/chk_adhoc_sched/check_adhoc_sched.py > /sso/sfw/LoadMgr/src/custom/chk_adhoc_sched/check_adhoc_sched.log 2>&1
30 2 * * 5 /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh > /sso/sfw/LoadMgr/src/custom/TestDataCheck/check_test_data.sh.log
*/5 7-17 * * * /sso/sfw/LoadMgr/custom/long_db_query_mon/long_db_query_mon.sh > /sso/sfw/LoadMgr/custom/long_db_query_mon/long_db_query_mon.log 2>&1
59 23 * * * source /jcp/projects/default/bin/crontab_backup.sh 1>~/crontab_backup.log 2>&1
39 22 06 06 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml --resume
46 20 02 07 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_thurs.xml --resume
27 22 11 07 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml --resume
33 22 15 08 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml --resume
14 19 24 08 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_mon.xml
04 21 30 08 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD_resume.json --schedule=/jcp/projects/hpmdo/conf/sched_sun_resume.xml
34 22 05 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml --resume
36 19 06 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sun.xml --skip_dependencies
03 18 20 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml
42 11 21 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sun.xml --skip_dependencies
05 16 25 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_fri_backups_psp2prod_sync.xml
18 22 26 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sat.xml --resume
00 22 27 09 * /sso/sfw/ghusps-toolkit/toolkit_venv/bin/python /sso/sfw/ghusps-toolkit/ams-toolkit/src/Toolkit/Controllers/ams_schedule_launcher.py --config_file=/jcp/projects/hpmdo/conf/JCP_MI_PROD.json --schedule=/jcp/projects/hpmdo/conf/sched_sun.xml --resume --skip_dependencies
